# rrr.conf
[program:rrr]
command=python3 /home/curator2/rrr/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/rrr/test.err.log
stdout_logfile=/home/curator2/rrr/test.out.log